# services/openrc
#!/sbin/openrc-run

name="shac"
command=

depend() {
    need docker
}

reload() {
    ebegin "Reloading configuration"
    $command
    eend $?
}